# ansible/roles/search_index/templates/rabbitmq-consumer.service
[Unit]
Description=RabbitMQ consumer

[Service]
Type=simple
ExecStart={{ etsin_reindexer_path }}/rabbitmq/rabbitmq_client.py

[Install]
WantedBy=multi-user.target
